# supervisord.conf
[supervisord]
nodaemon = true

[program:startup]
command = /scripts/startup.sh
autorestart = false

[program:mysql]
command = /usr/bin/mysqld_safe
stdout_logfile = /home/work/open-falcon/logs/mysql.log
stderr_logfile = /home/work/open-falcon/logs/mysql.err
autorestart=true

[program:redis]
user = redis
command = /usr/bin/redis-server /etc/redis/redis.conf
stdout_logfile = /home/work/open-falcon/logs/redis.log
stderr_logfile = /home/work/open-falcon/logs/redis.err

[program:alarm]
user = work
directory = /home/work/open-falcon/
command = %(program_name)/falcon-%(program_name) -c conf/%(program_name).cfg
stdout_logfile = logs/%(program_name).log
stderr_logfile = logs/%(program_name).err

[program:fe]
user = work
directory = /home/work/open-falcon/
command = %(program_name)/falcon-%(program_name) -c conf/%(program_name).cfg
stdout_logfile = logs/%(program_name).log
stderr_logfile = logs/%(program_name).err

[program:graph]
user = work
directory = /home/work/open-falcon/
command = %(program_name)/falcon-%(program_name) -c conf/%(program_name).cfg
stdout_logfile = logs/%(program_name).log
stderr_logfile = logs/%(program_name).err

[program:hbs]
user = work
directory = /home/work/open-falcon/
command = %(program_name)/falcon-%(program_name) -c conf/%(program_name).cfg
stdout_logfile = logs/%(program_name).log
stderr_logfile = logs/%(program_name).err

[program:judge]
user = work
directory = /home/work/open-falcon/
command = %(program_name)/falcon-%(program_name) -c conf/%(program_name).cfg
stdout_logfile = logs/%(program_name).log
stderr_logfile = logs/%(program_name).err

[program:query]
user = work
directory = /home/work/open-falcon/
command = %(program_name)/falcon-%(program_name) -c conf/%(program_name).cfg
stdout_logfile = logs/%(program_name).log
stderr_logfile = logs/%(program_name).err

[program:sender]
user = work
directory = /home/work/open-falcon/
command = %(program_name)/falcon-%(program_name) -c conf/%(program_name).cfg
stdout_logfile = logs/%(program_name).log
stderr_logfile = logs/%(program_name).err

[program:transfer]
user = work
directory = /home/work/open-falcon/
command = %(program_name)/falcon-%(program_name) -c conf/%(program_name).cfg
stdout_logfile = logs/%(program_name).log
stderr_logfile = logs/%(program_name).err

[program:dashboard]
user = work
directory = /home/work/open-falcon/dashboard
command = gunicorn -c gunicorn.conf wsgi:app -t 6000 
stdout_logfile = logs/%(program_name).log
stderr_logfile = logs/%(program_name).err

[program_name:portal]
user = work
directory = /home/work/open-falcon/portal
command = gunicorn -c gunicorn.conf wsgi:app -t 6000
stdout_logfile = logs/%(program_name).log
stderr_logfile = logs/%(program_name).err
