# supervisord.conf
[supervisord]
childlogdir = /home/work/open-falcon/logs/
logfile = /home/work/open-falcon/logs/supervisord.log
pidfile = /var/run/supervisord.pid
nodaemon = true
directory = /home/work/open-falcon

[unix_http_server]
file = /var/run/supervisor.sock
chmod = 0777
chown= nobody:nogroup

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///var/run/supervisor.sock

[program:startup]
command = bash /scripts/startup.sh
autorestart = false
priority = 1

[program:mysql]
command = /usr/bin/pidproxy /var/run/mysql.pid /usr/bin/mysqld_safe --defaults-file=/home/work/open-falcon/conf/my.cnf --skip-syslog
stdout_logfile = /home/work/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/work/open-falcon/logs/%(program_name)s.err
autostart=false
priority = 100

[program:redis]
user = redis
command = /usr/bin/redis-server /home/work/open-falcon/conf/redis.conf
stdout_logfile = /home/work/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/work/open-falcon/logs/%(program_name)s.err
autostart=false
priority = 100

[program:alarm]
user = work
directory = /home/work/open-falcon/%(program_name)s
command = /home/work/open-falcon/%(program_name)s/falcon-%(program_name)s -c /home/work/open-falcon/conf/%(program_name)s.cfg
stdout_logfile = /home/work/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/work/open-falcon/logs/%(program_name)s.err
autostart=false

[program:fe]
user = work
directory = /home/work/open-falcon/%(program_name)s
command = /home/work/open-falcon/%(program_name)s/falcon-%(program_name)s -c /home/work/open-falcon/conf/%(program_name)s.cfg
stdout_logfile = /home/work/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/work/open-falcon/logs/%(program_name)s.err
autostart=false

[program:graph]
user = work
directory = /home/work/open-falcon/%(program_name)s
command = /home/work/open-falcon/%(program_name)s/falcon-%(program_name)s -c /home/work/open-falcon/conf/%(program_name)s.cfg
stdout_logfile = /home/work/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/work/open-falcon/logs/%(program_name)s.err
autostart=false

[program:hbs]
user = work
directory = /home/work/open-falcon/%(program_name)s
command = /home/work/open-falcon/%(program_name)s/falcon-%(program_name)s -c /home/work/open-falcon/conf/%(program_name)s.cfg
stdout_logfile = /home/work/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/work/open-falcon/logs/%(program_name)s.err
autostart=false

[program:judge]
user = work
directory = /home/work/open-falcon/%(program_name)s
command = /home/work/open-falcon/%(program_name)s/falcon-%(program_name)s -c /home/work/open-falcon/conf/%(program_name)s.cfg
stdout_logfile = /home/work/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/work/open-falcon/logs/%(program_name)s.err
autostart=false

[program:query]
user = work
directory = /home/work/open-falcon/%(program_name)s
command = /home/work/open-falcon/%(program_name)s/falcon-%(program_name)s -c /home/work/open-falcon/conf/%(program_name)s.cfg
stdout_logfile = /home/work/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/work/open-falcon/logs/%(program_name)s.err
autostart=false

[program:sender]
user = work
directory = /home/work/open-falcon/%(program_name)s
command = /home/work/open-falcon/%(program_name)s/falcon-%(program_name)s -c /home/work/open-falcon/conf/%(program_name)s.cfg
stdout_logfile = /home/work/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/work/open-falcon/logs/%(program_name)s.err
autostart=false

[program:transfer]
user = work
directory = /home/work/open-falcon/%(program_name)s
command = /home/work/open-falcon/%(program_name)s/falcon-%(program_name)s -c /home/work/open-falcon/conf/%(program_name)s.cfg
stdout_logfile = /home/work/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/work/open-falcon/logs/%(program_name)s.err
autostart=false

[program:dashboard]
user = work
directory = /home/work/open-falcon/%(program_name)s
command = /home/work/open-falcon/%(program_name)s/env/bin/gunicorn -c gunicorn.conf wsgi:app -t 6000 
stdout_logfile = /home/work/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/work/open-falcon/logs/%(program_name)s.err
autostart=false

[program:portal]
user = work
directory = /home/work/open-falcon/%(program_name)s
command = /home/work/open-falcon/%(program_name)s/env/bin/gunicorn -c gunicorn.conf wsgi:app -t 6000 
stdout_logfile = /home/work/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/work/open-falcon/logs/%(program_name)s.err
autostart=false

[program:links]
user = work
directory = /home/work/open-falcon/%(program_name)s
command = /home/work/open-falcon/%(program_name)s/env/bin/gunicorn -c gunicorn.conf wsgi:app -t 6000 
stdout_logfile = /home/work/open-falcon/logs/%(program_name)s.log
stderr_logfile = /home/work/open-falcon/logs/%(program_name)s.err
autostart=false
